# docker/8.3/work.conf
[program:queue]
command =/usr/bin/php /data/www/laravel/artisan queue:work
user = www-data
autostart = true
autorestart= true