# docker/8.2/cron.conf
[program:crond]
command = /usr/sbin/cron -f
user = root
autostart = true


#[program:rsyslogd]
#command= /sbin/rsyslogd -n
#numprocs=1
#autostart=true
#autorestart=true