# docker/8.2/work.dev.conf
[program:queue]
command =/usr/bin/php /data/www/laravel/artisan queue:listen
user = www-data
autostart = true
autorestart= true